# yaguara-systemd-units/app.service
[Unit]
Description=Yaguara application

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes
