# yaguara-systemd-units/shopify-patch.timer
[Unit]
Description=Run shopify patch every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
